# cron/sum_up.cron
*/10 * * * * cd /path/to/IdolSongVote/ && carton exec -- perl scripts/sum_up.pl
